# speedtest.cron
58,18,38 * * * * root sudo /usr/bin/speedtest-cli --simple > /var/log/munin/speedtest.out.new && mv /var/log/munin/speedtest.out.new /var/log/munin/speedtest.out
